# systemd-timers-service/backup-pgsql-hourly.timer
[Unit]
Description=backup PostgreSQL databases timer

[Timer]
OnCalendar=*-*-* *:31:00

[Install]
WantedBy=timers.target
